# deploy/boatrace-public.service
[Unit]
Description=BOTRACE allowlisted public HTTP origin
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/boatrace-atlas
ExecStart=/usr/bin/python3 %h/boatrace-atlas/tools/serve_public.py
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
